# bin/config/88-local.rules
SUBSYSTEMS=="usb", ACTION=="add", RUN+="/home/pi/usb-motion.bash ADDED"
SUBSYSTEMS=="usb", ACTION=="remove", RUN+="/home/pi/usb-motion.bash REMOVED"
