# bin/config/file_lister.conf
description "Upstart job to kick off file_lister.py script."
author "Bubblers"
start on runlevel [2345]
exec /usr/local/bin/file_lister.py
