# bin/config/media-discovery3.service
[Unit]
Description=Media Discovery Service
After=multi-user.target

[Service]
Type=idle
ExecStart=/bin/bash /home/pi/tree_helper.sh

[Install]
WantedBy=multi-user.target
