# bin/config/media-discovery2.service
[Unit]
Description=Media Discovery Service
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python /home/pi/file_lister_v2.py

[Install]
WantedBy=multi-user.target
