# bin/config/thumbs-gen.service
[Unit]
Description=Thumbnail Generation Service
After=multi-user.target

[Service]
PIDFile=/home/pi/thumbs-gen.pid
ExecStart=/bin/bash /home/pi/thumbs-gen.bash

[Install]
WantedBy=multi-user.target
